# backend/TaskTracker.service
[Unit]
Description=TaskTracker Backend API
After=network.target

[Service]
Type=exec
User=www-data
WorkingDirectory=/var/www/tasktracker/backend
ExecStart=/usr/bin/dotnet TaskTracker.dll
Restart=always
RestartSec=10
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000
Environment=DOTNET_ROOT=/usr/lib/dotnet
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/dotnet

[Install]
WantedBy=multi-user.target 